# docker/services.d/cron/run
#!/bin/sh

printf "[SUPERVISOR] Starting cron service ...\n"
exec crond -l 8 -f
