# docker/services.d/cron/run
#!/bin/sh

printf "[SUPERVISOR] Starting cron service ...\n"
exec cron -f
