# docker/services.d/amuled/run
#!/bin/sh

printf "[SUPERVISOR] Starting amuled service ...\n"
exec setpriv --reuid "${AMULE_UID}" --regid "${AMULE_GID}" --init-groups amuled -c "${AMULE_HOME}" -o
