# docker/services.d/amuled/run
#!/bin/sh
exec su-exec "${AMULE_USER}" amuled -c "${AMULE_HOME}" -o
